# Dockerfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mariadb]
command=/opt/mariadb/start.sh
